# worker/supervisord.conf
; =============================================================================
; supervisord.conf — System Windykacja Worker
; Uruchamia dwa procesy w jednym kontenerze:
;   1. ARQ Worker  — przetwarza zadania z kolejki Redis
;   2. Uvicorn API — REST API na porcie 8001
; =============================================================================

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
loglevel=info
pidfile=/tmp/supervisord.pid
user=root

[program:arq_worker]
command=python -m arq worker.main.WorkerSettings
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
stopwaitsecs=30
stdout_logfile=/app/logs/arq_worker_stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/arq_worker_stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app"
priority=10

[program:api_server]
command=uvicorn worker.api.app:app --host 0.0.0.0 --port 8001 --workers 1 --log-level info
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
stopwaitsecs=10
stdout_logfile=/app/logs/api_server_stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/api_server_stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app"
priority=20

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[eventlistener:process_monitor]
command=python /app/worker/core/process_monitor.py
events=PROCESS_STATE
buffer_size=100
stdout_logfile=/app/logs/process_monitor.log
stderr_logfile=/app/logs/process_monitor_err.log